# archiso-profile/airootfs/etc/systemd/system/licos-setup.service
[Unit]
Description=licOS live environment setup
After=local-fs.target
Before=getty@tty1.service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/licos-setup

[Install]
WantedBy=sysinit.target
